# ordiri/examples/k8s/files/vault/systemd/cert-renewer@vault.service.d/override.conf
[Unit]
ExecStartPost=/usr/bin/env sh -c "! systemctl --quiet is-active vault.service || systemctl reload vault.service"