# ordiri/examples/k8s/files/vault/systemd/cert-renewer@vault.service.d/override.conf
[Unit]
Wants=reload-vault.service
Before=reload-vault.service