# ordiri/examples/k8s/files/vault/systemd/cert-renewer@vault.service.d/override.conf
[Service]
ExecStartPost=/usr/bin/env sh -c "systemctl reload vault.service"