# ordiri/examples/k8s/files/vault/systemd/vault.service.d/override.conf
[Unit]
Requires=network-online.target fetch-unseal-token.service
After=network-online.target fetch-unseal-token.service cert-renewer@vault.service