# ordiri/examples/k8s/files/vault/systemd/fetch-unseal-token.service
[Unit]
Description=Get the unseal token
Wants=network-online.target
Before=vault.service

[Service]
ExecStart=/bin/bash -c /sbin/fetch-unseal-token.sh

[Install]
RequiredBy=vault.service