# ordiri/examples/k8s/files/vault/systemd/reload-vault.service
[Unit]
Description=Reload vault
After=vault.service
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl reload vault.service