# ordiri/examples/k8s/files/kube-masters/systemd/kube-master-bootstrap.service
[Unit]
Description=Configure Kube Masters

[Service]
Type=oneshot
ConditionPathExists=!/etc/kubernetes/admin.conf
ExecStart=/bin/bash -c /sbin/bootstrap.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target