# ordiri/examples/k8s/files/vault-root/systemd/vault.service.d/override.conf
[Unit]
After=vault-root.mount
Wants=vault-root.mount